# install/cron/hivetool.cron
#Check Sensors every 5 Minutes
*/5 * * * * /home/HiveControl/scripts/system/currconditions.sh > /dev/null 2>&1

#Upload data once every 30 minutes to yardcontroller - disabled by default
#*/32 * * * * /home/HiveControl/scripts/data/synchive.sh > /dev/null 2>&1

#Take a Pic
#*/5 * * * * /home/HiveControl/scripts/image/takepic.sh > /dev/null 2>&1

#Get GDD Data Daily
00 12 * * * /home/HiveControl/scripts/weather/gdd/get_daily_insect_gdd.sh > /dev/null 2>&1

#Check for HiveControl Updates
*/32 * * * * /home/HiveControl/scripts/system/checkupgrades.sh > /dev/null 2>&1

#Get Pollen Count Daily
00 12 * * * /home/HiveControl/scripts/weather/pollen/getpollen.sh > /dev/null 2>&1

#Check Message queue once a minute
*/1 * * * * /home/HiveControl/scripts/system/msgqueue.sh > /dev/null 2>&1

#Get System Status
*/5 * * * * timeout -s KILL 1m /home/HiveControl/scripts/pimonitor/server_status.sh > /home/HiveControl/www/public_html/admin/server_status.txt

#Restart Our Pi GPIO Library Deamon on restart
@reboot              /usr/local/bin/pigpiod